# infra/pgbackrest/pgbackrest-full-backup.service
[Unit]
Description=pgBackRest full backup
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/bin/pgbackrest --stanza=mcpgov --type=full backup
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
